# vps-monitor/vps-monitor.service
[Unit]
Description=VPS Internet Monitor
After=network-online.target tailscaled.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vps-monitor
ExecStart=/opt/vps-monitor/venv/bin/python /opt/vps-monitor/monitor.py
Restart=always
RestartSec=10

# Environment
Environment=CONFIG_PATH=/opt/vps-monitor/config.yaml
# REQUIRED: Set your ntfy topic
Environment=NTFY_TOPIC=your-ntfy-topic-here
# Uncomment to override other settings
#Environment=LISTEN_PORT=5000
#Environment=HEARTBEAT_TIMEOUT=180

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-monitor

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vps-monitor
PrivateTmp=true

[Install]
WantedBy=multi-user.target
